# services-config/jd-controller-api.service
[Unit]
Description=JDownloader Controller API Server
After=network-online.target jdownloader.service
Wants=network-online.target

[Service]
Type=simple
User=ght
Group=ght
WorkingDirectory=/opt/jd2-controller/current
Environment="PATH=/opt/jd2-controller/current/venv/bin"
EnvironmentFile=/opt/jd2-controller/current/.env

# Start API server
ExecStart=/opt/jd2-controller/current/venv/bin/python /opt/jd2-controller/current/main.py api --prod

# Restart policy
Restart=on-failure
RestartSec=5

# Resource limits
Nice=5

# Logs
StandardOutput=journal
StandardError=journal
SyslogIdentifier=jd-controller-api

[Install]
WantedBy=multi-user.target
